# docker-files/supervisord.conf
[supervisord]
nodaemon=true

[program:artbot]
command=python /usr/src/app/artbot.py
redirect_stderr=true
stdout_logfile=/var/log/supervisor/artbot.log